# deploy/cron.d
* * * * * root python src/main.py check_user_status